# com.playwright.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.playwright.server</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>__NODE_PATH__</string>
        <string>__PROJECT_DIR__/server.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/logs/server.log</string>
    
    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/logs/server.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>DISPLAY</key>
        <string>:0</string>
    </dict>
</dict>
</plist>
